# data/etc/systemd/system/systemd-networkd-wait-online.service.d/waitany.conf
[Service]
ExecStart=
ExecStart=/lib/systemd/systemd-networkd-wait-online --any -i lo:missing -i eth0:missing --timeout=10 -q